# deploy/systemd/gigapdf-celery-billing.service
[Unit]
Description=GigaPDF Celery Billing Worker
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/gigapdf
Environment="PATH=/opt/gigapdf/.venv/bin"
EnvironmentFile=/opt/gigapdf/.env
ExecStart=/opt/gigapdf/.venv/bin/celery -A app.tasks.celery_app worker --loglevel=info --queues=billing --concurrency=2
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/var/log/gigapdf/celery-billing.log
StandardError=append:/var/log/gigapdf/celery-billing-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/gigapdf /var/lib/gigapdf /var/log/gigapdf

[Install]
WantedBy=multi-user.target
